# bitcoind/container-files/etc/service/bitcoind/run
#!/bin/bash

exec /usr/bin/bitcoind \
  -datadir=/root/.bitcoin \
  -rpcuser=kek \
  -rpcpassword=kek \
  -printtoconsole
